# 8-Pwntools-lab/Docker/root/etc/services.d/server/run
#!/usr/bin/with-contenv sh

/usr/bin/with-contenv /app/server
echo "starting self-compiled server"
exec /entry.sh
