# 10-Patching-lab/Docker/root/etc/services.d/server/run
#!/usr/bin/with-contenv sh

# echo "starting self-compiled server"
# exec /entry.sh
